# com.nativehp.p1102-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nativehp.p1102-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>PATH_TO_DRIVER_DIR/venv/bin/python3</string>
        <string>PATH_TO_DRIVER_DIR/p1102_daemon.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>PATH_TO_DRIVER_DIR/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>PATH_TO_DRIVER_DIR/daemon_error.log</string>
    <key>WorkingDirectory</key>
    <string>PATH_TO_DRIVER_DIR</string>
</dict>
</plist>
